# containers/edgemanage/supervisord.conf
[supervisord]
nodaemon=true

[program:edgemanage]
command=/usr/local/bin/edge_manage --daemonise --dnet dnext1
